# bluefield/misc/site-controller/otel-agent.service
[Unit]
Description=DPU OpenTelemetry mTLS certificate renewal agent (otel-agent kubelet static pod)
After=networking.service network-online.target time-sync.target containerd@mgmt.service
StartLimitIntervalSec=0
StartLimitBurst=0

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/share/otelcol-contrib/docker/otel-agent/otel-agent.sh
ExecStop=/bin/sh -c 'rm -f /etc/kubelet.d/otel-agent.yaml'
ExecStopPost=/bin/sh -c 'for i in $(seq 1 30); do \
    if ! /usr/bin/crictl ps --name otel-agent | grep -q otel-agent; then \
      exit 0; \
    fi; \
    sleep 2; \
  done; \
  echo "otel-agent container still running after timeout" >&2; \
  exit 1'
Restart=on-failure
RestartSec=30

User=root
Group=root

[Install]
WantedBy=multi-user.target
