# bluefield/misc/site-controller/DEBIAN/postinst
#!/bin/sh
set -e

. /usr/share/debconf/confmodule

if [ "$1" = "configure" ]; then
	# Ensure the core file directory exists
	# Must match /proc/sys/kernel/core_pattern
	mkdir -p /var/support/core

	# otel_post_install.sh is written to be harmless if run a second time.
	if ! /usr/local/sbin/otel_post_install.sh; then
		echo "otel_post_install.sh failed, will retry on next configure." >&2
	fi

	systemctl daemon-reload
	systemctl enable otelcol-contrib.service
	systemctl enable otel-agent.service

	# Do not depend on starting services now, as this package is installed early in setup
	# and should not hang or cause out of order steps in discovery workflows. The discovery
	# workflow will reboot shortly after this package is installed, so it is best to rely
	# on the reboot process to start the services enabled above.

	if [ -n "$2" ]; then
 		# The package is being upgraded and not installed for the first time, so
 		# this is not the discovery workflow and otelcol-contrib can restart without
 		# interfering. Restart is needed to pick up any configuration changes.

		systemctl restart otel-agent.service
		systemctl restart otelcol-contrib.service
	fi
fi
